# ops/cron/safe-zone-production.cron.example
# Safe Zone production cron examples.
# Install with: crontab ops/cron/safe-zone-production.cron.example

*/5 * * * * cd /opt/safe-zone && . ./.env && scripts/duckdns-update.sh >> logs/duckdns.log 2>&1
30 * * * * cd /opt/safe-zone && . ./.env && scripts/safe-zone.sh feed-sync >> logs/feed-sync.log 2>&1
15 3 * * * cd /opt/safe-zone && . ./.env && scripts/safe-zone.sh backup >> logs/backup.log 2>&1
